# etc_skel/xubuntu-14.04/cube_skel/.config/upstart/desktopOpen.conf
description "Desktop Open Task"
start on desktop-start
task
script
  # Disable screen blanking
  #echo "Desktop Open Task" >> ${HOME}/xxx_open_desktop
  xset s 0 0
end script
